# templates/supervisor_template.conf
[program:PROJECT_NAME]
command = SCRIPTS_DIR/gunicorn_exec.sh
user = SERVER_USER
stdout_logfile = RUN_DIR/gunicorn.log
redirect_stderr = true
autorestart = true
